# deploy/zerodte-shadow.service
[Unit]
# 0DTE GEX shadow pipeline + optional Grok 4.5 autonomous paper track.
# No live-order adapter exists; every fill remains inside PaperBroker.
Description=0DTE GEX shadow pipeline with Grok paper trader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zerodte
Group=zerodte
WorkingDirectory=/opt/zerodte

# Secrets live ONLY here (chmod 600, root:zerodte). Never in the repo or unit.
EnvironmentFile=/etc/zerodte/zerodte.env

# The journal, paper DB, Grok audit/cost ledger, and restart-safe Grok position
# snapshot persist under /var/lib/zerodte.
ExecStart=/opt/zerodte/venv/bin/python /opt/zerodte/grok_runner.py \
    --db /var/lib/zerodte/shadow.db \
    --paper-db /var/lib/zerodte/paper.sqlite \
    --live-state /var/lib/zerodte/live_state.json \
    --champion /var/lib/zerodte/configs/champion.json \
    --interval 60 \
    --no-legacy-directional-tilt

# Auto-restart on crash; the loop self-handles market hours, so this only
# fires on real faults (network, API, bug). RestartSec avoids hot-looping.
Restart=on-failure
RestartSec=10

# --- Hardening: the box holds API keys, so lock the process down. ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
# Only the journal directory is writable; everything else is read-only.
ReadWritePaths=/var/lib/zerodte

StandardOutput=journal
StandardError=journal
SyslogIdentifier=zerodte-shadow

[Install]
WantedBy=multi-user.target
